# scripts/auto-logout-key/85-yubikey.rules
ACTION=="remove", SUBSYSTEM=="usb", ENV{PRODUCT}=="1050/402/543", RUN+="/usr/bin/loginctl lock-sessions"
